# roles/harbor/files/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=/usr/local/php/sbin/php-fpm 

[program:nginx]
command=/usr/local/nginx/sbin/nginx